# extra/twt.service.template
[Unit]
Description=twt
PartOf=graphical-session.target

[Service]
Environment=DISPLAY=:0
ExecStart=PATH_TO_EXECUTABLE run
Restart=always

[Install]
WantedBy=default.target
